# scripts/schedulers/com.xgh.models.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS Apple Silicon only — for Linux/Intel Mac use Ollama (systemctl --user status xgh-models) -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.xgh.models</string>
  <key>ProgramArguments</key>
  <array>
    <string>XGH_VLLM_BIN</string>
    <string>serve</string>
    <string>XGH_LLM_MODEL</string>
    <string>--embedding-model</string>
    <string>XGH_EMBED_MODEL</string>
    <string>--port</string>
    <string>XGH_MODEL_PORT</string>
    <string>--host</string>
    <string>XGH_MODEL_HOST_PLACEHOLDER</string>
  </array>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>XGH_LOG_DIR/vllm-mlx.log</string>
  <key>StandardErrorPath</key>
  <string>XGH_LOG_DIR/vllm-mlx.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>XGH_USER_HOME</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
